# tarsnap.timer
[Unit]
Description=Run Tarsnap weekly

[Timer]
OnCalendar=weekly

[Install]
WantedBy=tarsnap.target